# deploy/launchd/com.vanta.wake.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  VANTA always-on wake listener (template).
  Placeholders __PYTHON__ / __WORKDIR__ / __LOG__ are substituted by
  scripts/install_wake_listener.sh at install time (launchd needs absolute paths).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vanta.wake.listener</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>openjarvis.speech.wake_listener</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>__WORKDIR__/src</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__LOG__</string>

  <key>StandardErrorPath</key>
  <string>__LOG__</string>
</dict>
</plist>
